# profile.d/activeqmlimportpath.sh.cmake
 
export QML_IMPORT_PATH=@KDE_INSTALL_PREFIX@/lib/kde4/platformimports/touch:@KDE_INSTALL_PREFIX@/lib/kde4/imports
