# profile.d/disablecursors.sh.cmake
 
XCURSOR_PATH=@KDE_INSTALL_PREFIX@/share/icons::~/.icons:/usr/share/icons:/usr/share/pixmaps:/usr/X11R6/lib/X11/icons
XCURSOR_THEME=plasmamobilemouse
