# skel/dotconfig/autostart/korgac.desktop
[Desktop Entry]
Exec=korgac
Hidden=true
OnlyShowIn=KDE;
X-KDE-autostart-phase=BaseDesktop
